# resources/upstart/rest.conf
description "rest"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid rest
setgid rest
chdir /home/rest/generic-rest-service

exec /home/rest/generic-rest-service/bin/production.sh
